# ws-server/websocket-server.service
[Unit]
Description=websocket-server

[Service]
ExecStart= /home/ec2-user/.nvm/versions/node/v22.12.0/bin/node --max-old-space-size=2048 /data/prsm/node_modules/.bin/y-websocket-server
Restart=always
RestartSec=10
User=root
Group=root
Environment='YPERSISTENCE=/data/prsm/dbDir'
# Memory limit to prevent system exhaustion
MemoryMax=2.5G
MemoryAccounting=true
SyslogIdentifier=websocket-server

[Install]
WantedBy=multi-user.target
